# cookbooks/blog/files/etc/cron.d/blog
@reboot webadm cp -pr /home/webadm/works/public/* /var/www/blog/
