# cookbooks/everun/files/etc/cron.d/everun-blog
@reboot webadm cp -pr /home/webadm/works/everun/* /var/www/everun/
