# cookbooks/docker/files/etc/cron.d/docker-housekeep
MAILTO=""

DOCKER=/usr/bin/docker

0 13 * * *    root ${DOCKER} image prune -f
0 23 * * *    root ${DOCKER} container prune -f
