# cookbooks/wekan/files/etc/supervisor/conf.d/wekan.conf
[program:wekan]
command=/home/wekan/app/wekan.sh
stdout_logfile=/var/log/supervisor/wekan.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
autorestart=true
user=wekan
