# cookbooks/digdag/files/etc/supervisor/conf.d/digdag.conf
[program:digdag]
command=/etc/digdag/digdag.sh
stdout_logfile=/var/log/supervisor/digdag.log
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
autorestart=true
