# cookbooks/nginx/files/etc/cron.d/ssl
2 2 */14 * * root /home/webadm/bin/ssl_renewal.sh

