# cookbooks/prometheus/files/lib/systemd/system/promtail-prometheus.service
[Unit]
Description=Grafana Promtail for Prometheus
Documentation=https://github.com/grafana/loki
After=network-online.target

[Service]
User=root
Restart=always
ExecStart=/usr/local/bin/promtail --config.file=/etc/promtail/prometheus.yaml

[Install]
WantedBy=multi-user.target
