# cookbooks/statup/files/etc/cron.d/statup
12 3 * * * root openssl rand 48 > /etc/letsencrypt/live/status.kazu634.com/ticket.key
