# cookbooks/loki/files/etc/systemd/system/promtail-loki.service
[Unit]
Description=Grafana Promtail
Documentation=https://github.com/grafana/loki
After=network-online.target

[Service]
User=root
Restart=always
ExecStart=/usr/local/bin/promtail --config.file=/etc/promtail/loki.yaml

[Install]
WantedBy=multi-user.target
